# crontab
59 * * * * if [[ $(yy --alive) == hived ]]; then :;else yy --demon &;fi
