# services/kmonad/run
#!/bin/sh

exec /usr/bin/kmonad /home/devi/scripts/keymap.kbd 2>&1
