# services/clipd/run
#!/bin/sh
exec \
  env \
  XDG_RUNTIME_DIR="/run/user/1000" \
  DISPLAY=:0 \
  /usr/bin/lua5.3 \
  /home/devi/devi/lclip.git/in_memory/lclipd.lua \
  > /dev/null 2>&1
