# services/clipd/run
#!/bin/sh
exec 2>&1
DISPLAY=:0.0 exec /home/devi/scripts/bin/clipd
